# claude/ec2-agents/LaunchAgents/com.mai.claude-gchat-default.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mai.claude-gchat-default</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ec2-user/agents/start-claude-default.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ec2-user/agents/default</string>
    <key>StandardOutPath</key>
    <string>/tmp/com.mai.claude-gchat-default.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.mai.claude-gchat-default.err</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
